# sh/crontab
*/3 * * * *     /usr/sbin/logrotate -f /etc/logrotate.d/verenav.conf
*/10 * * * *    nice -n 15 /upload.sh
