# sh/crontab
*/4 * * * *     /usr/sbin/logrotate -f /etc/logrotate.d/verenav.conf
*/15 * * * *    nice -n 15 /upload.sh
